# build/config/systemd/mount_spade_volumes.service
# Systemd service at /etc/systemd/system/mount_spade_volumes.service

[Unit]
After=filesystem

[Service]
Type=oneshot
ExecStartPre=/bin/ln -sFTf /mnt /opt/science/spade/data
ExecStartPre=/usr/sbin/ebslvm vgebs lvebs /mnt
RemainAfterExit=true
