# build/config/upstart/mount_spade_volumes.conf
# upstart file to mount volume

start on filesystem

emits ebslvm-mounted

script
  ln -sFTf /mnt /opt/science/spade/data
  /usr/sbin/ebslvm vgebs lvebs /mnt
  initctl emit ebslvm-mounted
end script
