# build/config/upstart/spade.conf
# upstart file to run spade

start on ebslvm-mounted and stopped cloud-final
stop on starting rc RUNLEVEL=[016]

emits spade-running

respawn
respawn limit 10 5

kill timeout 600
kill signal SIGINT

limit nofile 4096 4096

script
  exec /opt/science/spade/bin/run_spade.sh
  initctl emit spade-running
end script
