# build/config/mount_spade_volumes.conf
# upstart file to mount volume

start on runlevel [2345]
stop on runlevel [016]

script
  /usr/sbin/ebslvm vgebs lvebs /mnt
  ln -sFTf /mnt /opt/science/spade/data
  emit mounted
end script
